# os-mosdns/pkg-plist
bin/mosdns
etc/inc/plugins.inc.d/mosdns.inc
etc/rc.d/mosdns
@sample etc/mosdns/config.yaml.sample
opnsense/mvc/app/controllers/OPNsense/MosDNS/AdvancedController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/AdvancedController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/GeneralController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/LogsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/PluginsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/ServiceController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/SettingsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/Api/StatisticsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/IndexController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/LogsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/PluginsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/StatisticsController.php
opnsense/mvc/app/controllers/OPNsense/MosDNS/forms/advanced.xml
opnsense/mvc/app/controllers/OPNsense/MosDNS/forms/general.xml
opnsense/mvc/app/controllers/OPNsense/MosDNS/forms/plugins.xml
opnsense/mvc/app/controllers/OPNsense/MosDNS/forms/upstreams.xml
opnsense/mvc/app/models/OPNsense/MosDNS/ACL/ACL.xml
opnsense/mvc/app/models/OPNsense/MosDNS/Menu/Menu.xml
opnsense/mvc/app/models/OPNsense/MosDNS/MosDNS.php
opnsense/mvc/app/models/OPNsense/MosDNS/MosDNS.xml
opnsense/mvc/app/views/OPNsense/MosDNS/advanced.volt
opnsense/mvc/app/views/OPNsense/MosDNS/general.volt
opnsense/mvc/app/views/OPNsense/MosDNS/index.volt
opnsense/mvc/app/views/OPNsense/MosDNS/logs.volt
opnsense/mvc/app/views/OPNsense/MosDNS/plugins.volt
opnsense/mvc/app/views/OPNsense/MosDNS/statistics.volt
opnsense/scripts/OPNsense/MosDNS/generate_config.sh
opnsense/scripts/OPNsense/MosDNS/setup.sh
opnsense/service/conf/actions.d/actions_mosdns.conf
opnsense/service/templates/OPNsense/MosDNS/mosdns.conf
@dir etc/mosdns
@dir opnsense/mvc/app/controllers/OPNsense/MosDNS
@dir opnsense/mvc/app/models/OPNsense/MosDNS
@dir opnsense/mvc/app/views/OPNsense/MosDNS
@dir opnsense/scripts/OPNsense/MosDNS
@dir opnsense/service/templates/OPNsense/MosDNS